# autostart/startup.sh.desktop
[Desktop Entry]
Exec=/home/jdyer/bin/startup.sh
Icon=dialog-scripts
Name=startup.sh
Path=
Type=Application
X-KDE-AutostartScript=true
